# deploy/hive.service
[Unit]
Description=HIVE.AI - Reine Nu v0.3.0 - HIVE.WORK
After=network.target

[Service]
User=hive
Group=hive
WorkingDirectory=/opt/hive
Environment="PATH=/opt/hive/venv/bin"
EnvironmentFile=/opt/hive/.env
ExecStart=/opt/hive/venv/bin/gunicorn wsgi:app \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --timeout 120 \
    --access-logfile /var/log/hive/access.log \
    --error-logfile /var/log/hive/error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
